# fs/etc/profile.d/more-modern-aliases.sh
alias .h='cd $(zoxide query --interactive)'
alias cat='bat --style plain'
alias ls='eza --icons --color=always'
alias lt='ls -l --tree --color=always'
