# fs/etc/profile.d/ngrok-path.sh
if [[ "$PATH" != *"/opt/ngrok/bin"* ]]; then
  export PATH="/opt/ngrok/bin:$PATH"
fi
